# packages/rpm/falcon-object.service
[Unit]
Description=Falcon - Object Server
Wants=network-online.target
After=syslog.target network-online.target

[Service]
Type=simple
LimitNOFILE=102400
User=swift
Restart=always
ExecStart=/usr/local/bin/falcon object -c /etc/swift/object-server.conf
TimeoutStopSec=60
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
